# deploy/services/ngrok.service
[Unit]
Description=Ngrok public URL service
After=frontend.service

[Service]
User=andreierdna
WorkingDirectory=/home/andreierdna/projects/mainframe/backend
ExecStart=ngrok http 7777
ExecStartPost=/home/andreierdna/.virtualenvs/mainframe/bin/python manage.py set_hooks
Restart=always

[Install]
WantedBy=multi-user.target
