# deploy/services/wifi-repair.service
[Unit]
Description=Wi-Fi auto-repair if disconnected
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/rpi/projects/mainframe/deploy/wifi-repair.sh

